# b1-pgsql-connection-cleaner.cron
0 */2 * * * root /usr/bin/b1-pgsql-connection-cleaner
